# .config/systemd/user/dotfiles-snapshot.service
[Unit]
Description=Snapshot installed packages to dotfiles repo

[Service]
Type=oneshot
ExecStart=%h/project/dotfiles/scripts/snapshot-packages.sh
